# templates/daemon/linux.service
[Unit]
Description=OpenClaw Claude Session Keepalive
After=network.target

[Service]
Type=oneshot
ExecStart=/bin/bash {{SCRIPT_PATH}}
Restart=on-failure
RestartSec=30

[Install]
WantedBy=default.target
